# config/supervisor/admin.conf
[program:admin_server]
command=gunicorn_django -c /cube/app/admin/server/Adminboard/gunicorn.conf
directory=/cube/app/admin/server/Adminboard/
user=www-data
autostart=true
autorestart=true
stdout_logfile=/cube/app/admin/log/admin_supervisor.log
loglevel=warn
