# config/supervisor/admin.conf
[program:admin_server]
command=/usr/bin/gunicorn_django -c /cube/app/admin/server/Adminboard/gunicorn.conf
directory=/cube/app/admin/server/Adminboard/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/admin/log/admin_server.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_logfile=/cube/app/admin/log/admin_server.error_log
stderr_logfile_maxbytes=50MB
loglevel=warn
