# config/supervisor/api.conf
[group:api_server]
programs=api_server-3001,api_server-3002,api_server-3003,api_server-3004,api_server-3005,api_server-3006,api_server-3007,api_server-3008,api_server-3009,api_server-3010

[program:api_server-3001]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3001
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3001.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3002]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3002
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3002.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3003]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3003
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3003.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3004]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3004
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3004.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn


[program:api_server-3005]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3005
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3005.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3006]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3006
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3006.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn


[program:api_server-3007]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3007
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3007.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3008]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3008
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3008.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn


[program:api_server-3009]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3009
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3009.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn

[program:api_server-3010]
command=/cube/bin/python /cube/app/api/server/api/server.py --port=3010
directory=/cube/app/api/server/api/
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/cube/app/api/log/api_server-3010.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=warn
