# config/supervisor/dashboard.conf
[program:dashboard]
command=/cube/bin/gunicorn_django -c /cube/app/dashboard/server/dashboard/gunicorn.conf
directory=/cube/app/dashboard/server/dashboard
user=www-data
autostart=true
autorestart=true
stdout_logfile=/cube/app/dashboard/log/supervisord_dashboard.log
redirect_stderr=true
